# cookbooks/chef/templates/default/chef-client.service.erb
[Unit]
Description=Chef Client
After=network.target

[Service]
ExecStart=/usr/bin/chef-client -i 1800 -s 20
Restart=on-failure

[Install]
WantedBy=multi-user.target
