# cookbooks/chef/templates/default/chef-client.conf.erb
# chef-client
#
# Startup script for chef-client

description "starts up chef-client in daemon mode"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

script
        exec /usr/bin/chef-client -i 1800 -s 20
end script

respawn
