# cookbooks/chef/templates/default/repository-backup.cron.erb
#!/bin/sh

T=$(mktemp -d -t -p /var/tmp chef-repository.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=chef-repository-$D.tar.gz

ln -s /var/lib/git/chef.git $T/chef-repository-$D

export GZIP="--rsyncable -9"

nice tar --create --gzip --dereference --directory=$T --file=$T/$B chef-repository-$D
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
