# cookbooks/chef/templates/default/server-backup.cron.erb
#!/bin/bash

# DO NOT EDIT - This file is being maintained by Chef
set -euo pipefail

export ZSTD_CLEVEL=16

T=$(mktemp -d -t -p /var/tmp cinc-server.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B="cinc-server-$D.tar.zst"

mkdir "$T/cinc-server-$D"
chgrp opscode-pgsql "$T" "$T/cinc-server-$D"
chmod g+rwx "$T" "$T/cinc-server-$D"
sudo -u opscode-pgsql /opt/cinc-project/embedded/bin/pg_dumpall --file="$T/cinc-server-$D/cinc.dmp" --clean
ln -s /var/opt/cinc-project/bookshelf/data "$T/cinc-server-$D/bookshelf"

nice tar --create --dereference --directory="$T" \
  --sort=name \
  --warning=no-file-removed \
  "cinc-server-$D" | nice zstd -T0 --quiet --long --rsyncable -o "$T/$B"

nice rsync --preallocate --fuzzy "$T/$B" backup.openstreetmap.org::backup

rm -rf "$T"
