# cookbooks/db/templates/default/backup.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

00 02 * * 1 osmbackup /usr/local/bin/backup-db
