# cookbooks/tools/templates/default/cron.service.erb
[Service]
Environment=TZ=<%= node[:timezone] -%>
