# cookbooks/imagery/templates/default/mapserv_fcgi.conf.erb
# DO NOT EDIT - This file is being maintained by Chef

description     "starts mapserver in daemon mode for fastcgi"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

respawn

limit nofile 16384 16384

# https://trac.osgeo.org/mapserver/wiki/EnvironmentVariables
env MS_MAPFILE="/srv/imagery/mapserver/layer-<%= @name %>.map"
env MS_MAP_PATTERN="^/srv/imagery/mapserver/"
env MS_ERRORFILE="/tmp/mapserver-layer-<%= @name %>.log"
env MS_DEBUGLEVEL="5"

pre-start script
    mkdir -p -m 0775 /var/run/mapserver-fastcgi
    chown imagery:imagery /var/run/mapserver-fastcgi
end script

script
    exec /usr/bin/spawn-fcgi \
        -n \
        -u imagery \
        -g imagery \
        -s /var/run/mapserver-fastcgi/layer-<%= @name %>.socket \
        -M 0666 \
        -f /usr/lib/cgi-bin/mapserv
end script
