# cookbooks/blog/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp osm-blog.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=osm-blog-$D.tar.gz

mkdir $T/osm-blog-$D
echo '[mysqldump]' > $T/mysqldump.opts
echo 'user=osm-blog-user' >> $T/mysqldump.opts
echo 'password=<%= @passwords["osm-blog-user"] %>' >> $T/mysqldump.opts
mysqldump --defaults-file=$T/mysqldump.opts --opt --no-tablespaces osm-blog > $T/osm-blog-$D/osm-blog.sql
ln -s /srv/blog.openstreetmap.org $T/osm-blog-$D/www

export RSYNC_RSH="ssh -ax"

nice tar --create --dereference --directory=$T --warning=no-file-changed osm-blog-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
