# cookbooks/blog/templates/default/backup-staging.cron.erb
#!/bin/bash

# DO NOT EDIT - This file is being maintained by Chef
set -euo pipefail

export ZSTD_CLEVEL=11
export ZSTD_NBTHREADS=0

T=$(mktemp -d -t -p /var/tmp osm-blog-staging.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B="osm-blog-staging-$D.tar.zst"

mkdir "$T/osm-blog-staging-$D"
echo '[mysqldump]' > "$T/mysqldump.opts"
echo 'user=osm-blog-staging-user' >> "$T/mysqldump.opts"
echo 'password=<%= @passwords["osm-blog-staging-user"] %>' >> "$T/mysqldump.opts"
mysqldump --defaults-file="$T/mysqldump.opts" --opt --no-tablespaces --max-allowed-packet=1G osm-blog-staging > "$T/osm-blog-staging-$D/osm-blog-staging.sql"
ln -s /srv/staging.blog.openstreetmap.org "$T/osm-blog-staging-$D/www"

set +e
nice tar --create --dereference --directory="$T" \
  --sort=name \
  --warning=no-file-changed \
  --warning=no-file-removed \
  --ignore-failed-read \
  "osm-blog-staging-$D" | nice zstd --quiet --long --rsyncable -o "$T/$B"
nice rsync --preallocate --fuzzy "$T/$B" backup.openstreetmap.org::backup

rm -rf "$T"
