# cookbooks/nginx/templates/nginx-old-cache-cleanup.cron.erb
# Cleanup old nginx cache storage
15 23 * * * www-data /usr/bin/timeout 6h /usr/local/bin/nginx-old-cache-cleanup
