# cookbooks/munin/templates/default/rrddump.cron.erb
MAILTO=zerebubuth@gmail.com
# do the dump & cleanup in the early hours of the morning
43 3 * * * www-data nice /usr/local/bin/rrddump
