# cookbooks/trac/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp trac.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=trac-$D.tar.gz

trac-admin /var/lib/trac hotcopy $T/trac-$D > /dev/null

export RSYNC_RSH="ssh -ax"

nice tar --create --directory=$T trac-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
