# cookbooks/stats/templates/default/repack.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

MAILTO=horntail-www-data-cron@firefishy.com
0 3 * * 7 www-data nice /usr/local/bin/repack-archived-logs
