# cookbooks/donate/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp osm-donate.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=osm-donate-$D.tar.gz

mkdir $T/osm-donate-$D
echo '[mysqldump]' > $T/mysqldump.opts
echo 'user=donate' >> $T/mysqldump.opts
echo 'password=<%= @passwords["database"] %>' >> $T/mysqldump.opts
mysqldump --defaults-file=$T/mysqldump.opts --opt --no-tablespaces donate > $T/osm-donate-$D/osm-donate.sql
ln -s /srv/donate.openstreetmap.org $T/osm-donate-$D/www

export RSYNC_RSH="ssh -ax"

nice tar --create --dereference --directory=$T osm-donate-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
