# cookbooks/mediawiki/templates/default/mediawiki-backup.cron.erb
#!/bin/bash
T=`mktemp -d -t -p /var/tmp mediawiki-<%= @name %>.XXXXXXXXXX`
D=`date +%Y-%m-%d`
B=wiki-<%= @name %>-$D.tar.gz

mkdir $T/wiki-<%= @name %>-$D
mysqldump --user="<%= @database_params[:username] %>" --password="<%= @database_params[:password] %>" --opt --skip-lock-tables --single-transaction "<%= @database_params[:name] %>" | xz -9c > $T/wiki-<%= @name %>-$D/wiki.sql.xz
ln -s <%= @directory %>  $T/wiki-<%= @name %>-$D/www
export GZIP="--rsyncable -9" #make backup rsyncable
nice tar --create --gzip --dereference --directory=$T --exclude=wiki-<%= @name %>-$D/w/images/thumb --exclude=wiki-<%= @name %>-$D/w/.git --file=$T/$B wiki-<%= @name %>-$D
nice rsync $T/$B backup::backup --fuzzy

rm -rf $T
