# cookbooks/nominatim/templates/default/nominatim.cron.erb
# Maintained by chef. DO NOT EDIT.

MAILTO=<%= @mailto %>

* * * * * lonvia <%= @bin_directory %>/cron_ipanalyse.py <%= node[:nominatim][:logdir] %>/nominatim.openstreetmap.org-access.log
<% if node[:nominatim][:enable_backup] -%>
00 3 1 * * nominatim /usr/local/bin/backup-nominatim
<% end -%>
20 0 * * * postgres /usr/local/bin/vacuum-db-nominatim
<% if node[:nominatim][:state] == "master" -%>
05 */4 * * * postgres /usr/local/bin/clean-db-nominatim
<% end -%>
18 1 * * * nominatim touch <%= @update_maintenance_trigger %>
