# cookbooks/community/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp community.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=community-$D.tar.gz

mkdir $T/community-$D
ln -s /srv/community.openstreetmap.org/docker/containers $T/community-$D/containers
ln -s /srv/community.openstreetmap.org/shared/web-only $T/community-$D/shared-web-only
ln -s /srv/community.openstreetmap.org/shared/data/redis_data $T/community-$D/shared-data-redis_data
ln -s /srv/community.openstreetmap.org/shared/data/postgres_backup $T/community-$D/shared-data-postgres_backup

export RSYNC_RSH="ssh -ax"

nice tar --create --numeric-owner --dereference --directory=$T community-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
