# cookbooks/prometheus/templates/default/backup.cron.erb
#!/bin/bash

# DO NOT EDIT - This file is being maintained by Chef
set -euo pipefail

export ZSTD_CLEVEL=11
export ZSTD_NBTHREADS=0

T=$(mktemp -d -t -p /var/tmp prometheus.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B="prometheus-$D.tar.zst"

mkdir "$T/prometheus-$D"

ln -s /var/lib/prometheus/alertmanager "$T/prometheus-$D/alertmanager"
ln -s /var/lib/grafana "$T/prometheus-$D/grafana"

nice tar --create --dereference --directory="$T" \
  --sort=name \
  "prometheus-$D" | nice zstd --quiet --long --rsyncable -o "$T/$B"

nice rsync --preallocate --fuzzy "$T/$B" backup.openstreetmap.org::backup

rm -rf "$T"
