# cookbooks/thinkup/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp thinkup.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=thinkup-$D.tar.gz

mkdir $T/thinkup-$D
mysqldump --user=thinkup --password=<%= @passwords["database"] %> --opt thinkup > $T/thinkup-$D/thinkup.sql
ln -s /srv/thinkup.openstreetmap.org $T/thinkup-$D/www

export GZIP="--rsyncable -9"

nice tar --create --gzip --dereference --directory=$T --file=$T/$B thinkup-$D
nice rsync --bwlimit=500 $T/$B backup::backup

rm -rf $T
