# cookbooks/hardware/templates/default/tty.conf.erb
# DO NOT EDIT - This file is being maintained by Chef

# ttyS<%= @unit %> - getty
#
# This service maintains a getty on ttyS<%= @unit %> from the point the system is
# started until it is shut down again.

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn
exec /sbin/getty -L <%= @speed %> ttyS<%= @unit %> vt102
