# cookbooks/web/templates/default/cleanup.cron.erb
#!/bin/sh

<%= node[:web][:base_directory] %>/rails/script/cleanup
