# cookbooks/web/templates/default/cleanup.cron.erb
#!/bin/sh

<%= @ruby %> <%= @directory %>/script/cleanup
