# cookbooks/web/templates/default/rails.cron.erb
#!/bin/sh

# Cleanup temporary directory
find <%= @directory %>/tmp -mtime +28 -print0 | xargs -0r rm
