# cookbooks/web/templates/default/rails.cron.erb
#!/bin/sh

# Cleanup temporary directory
if [ -d <%= @directory %>/tmp ]
then
    find <%= @directory %>/tmp \( \( -type f -mtime +28 \) -o \( -type d -empty \) \) -delete
fi

# Cleanup assets
if [ -x /usr/local/bin/cleanup-rails-assets ]
then
    /usr/local/bin/cleanup-rails-assets <%= @directory %>
fi
