# cookbooks/web/templates/default/rails.cron.erb
#!/bin/sh

# Cleanup temporary directory
find <%= @directory %>/tmp \( \( -type f -mtime +28 \) -o \( -type d -empty \) \) -delete
