# cookbooks/web/templates/default/statistics.cron.erb
# DO NOT EDIT - This file is being maintained by Chef
0 0 * * * rails /usr/local/bin/statistics
