# cookbooks/squid/templates/default/squid.service.erb
[Unit]
Description=Squid caching proxy
After=network.target nss-lookup.target

[Service]
LimitNOFILE=65536
Environment=SQUID_ARGS=-D
EnvironmentFile=/etc/default/squid
ExecStartPre=/usr/sbin/squid $SQUID_ARGS -z
ExecStart=/usr/sbin/squid -N $SQUID_ARGS
ExecReload=/usr/sbin/squid -k reconfigure
ExecStop=/usr/sbin/squid -k shutdown
Restart=on-failure
TimeoutSec=0

[Install]
WantedBy=multi-user.target
