# cookbooks/git/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp git.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=git-$D.tar.gz

ln -s /var/lib/git $T/git-$D

export RSYNC_RSH="ssh -ax"

nice tar --create --dereference --directory=$T git-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
