# cookbooks/git/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp git.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=<%= node[:git][:backup] %>-$D.tar.gz

ln -s /var/lib/git $T/git-$D

export GZIP="--rsyncable -9"
export RSYNC_RSH="ssh -ax"

nice tar --create --gzip --dereference --directory=$T --file=$T/$B git-$D
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
