# cookbooks/tile/templates/default/cleanup-tiles.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

<% @directories.each do |directory| -%>
0 1 * * * www-data ionice -c 3 /usr/local/bin/cleanup-tiles <%= directory %>
<% end -%>
