# cookbooks/tile/templates/default/export.cron.erb
#!/bin/sh

# Removes stale tmp files from the export tab
exec find /tmp -name 'export??????' -mmin +60 -delete
