# cookbooks/tile/templates/default/render-lowzoom.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

0 1 1-7 * * tile test $(date +\%u) -eq 7 && /usr/local/bin/render-lowzoom
