# cookbooks/planet/templates/default/old-planet-file-cleanup.cron.erb
# DO NOT EDIT - This file is being maintained by Chef
MAILTO=zerebubuth@gmail.com
# run this on the first monday of the month at 3:44am
44 3 1-7 * * www-data test $(date +\%u) -eq 1 && /usr/local/bin/old-planet-file-cleanup --dry-run --debug
