# cookbooks/planet/templates/default/planet-notes-dump.cron.erb
MAILTO=grant-smaug@firefishy.com
0 3 * * * www-data /usr/local/bin/planet-notes-dump
