# cookbooks/planet/templates/default/planet-update.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

17 1 * * * planet /usr/local/bin/planet-update
