# cookbooks/planet/templates/default/replication.cron.erb
# DO NOT EDIT - This file is being maintained by Chef

TZ=UTC

MAILTO=zerebubuth@gmail.com

0 7 * * * planet /usr/local/bin/users-agreed
0 17 * * * planet /usr/local/bin/users-deleted
* * * * * planet /usr/local/bin/replicate-changesets /etc/replication/changesets.conf

MAILTO=brett@bretth.com
LD_PRELOAD=/opt/flush/flush.so

* * * * * planet /usr/local/bin/osmosis -q --replicate-apidb authFile=/etc/replication/auth.conf validateSchemaVersion=false --write-replication workingDirectory=/store/planet/replication/minute
2,7,12,17 * * * * planet /home/bretth/bin/osmosis -q --merge-replication-files workingDirectory=/var/lib/replication/hour
5,10,15,20 * * * * planet /home/bretth/bin/osmosis -q --merge-replication-files workingDirectory=/var/lib/replication/day
