# cookbooks/subversion/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp svn.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=svn-$D.tar.gz

nice svnadmin hotcopy /var/lib/subversion/repos/openstreetmap $T/svn-$D > /dev/null

export GZIP="--rsyncable -9"
export RSYNC_RSH="ssh -ax"

nice tar --create --gzip --directory=$T --file=$T/$B svn-$D
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
