# cookbooks/subversion/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp svn.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=svn-$D.tar.gz

nice svnadmin hotcopy /var/lib/subversion/repos/openstreetmap $T/svn-$D > /dev/null

export RSYNC_RSH="ssh -ax"

nice tar --create --directory=$T svn-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
