# cookbooks/mailman/templates/default/backup.cron.erb
#!/bin/sh

# DO NOT EDIT - This file is being maintained by Chef

T=$(mktemp -d -t -p /var/tmp lists.XXXXXXXXXX)
D=$(date +%Y-%m-%d)
B=lists-$D.tar.gz

mkdir $T/lists-$D
ln -s /var/lib/mailman $T/lists-$D/mailman

export RSYNC_RSH="ssh -ax"

nice tar --create --dereference --warning=no-file-changed --directory=$T lists-$D | nice gzip --rsyncable -9 > $T/$B
nice rsync --preallocate --fuzzy $T/$B backup::backup

rm -rf $T
