# cookbooks/tilelog/templates/default/tilelog.cron.erb
MAILTO=zerebubuth@gmail.com
17 7 * * * www-data /usr/local/bin/tilelog
